# python-daemon/python_daemon.conf
[program:daemon]
command=/opt/python_daemon.py
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=unexpected
startsecs=1
startretries=30
stopsignal=QUIT
stopwaitsecs=10
